# scripts/dune-stack.service
# dune-stack.service - start the Dune Awakening stack on boot
#
# This is a TEMPLATE for packagers and manual installs.
# Prefer the portable generator, which detects your current user and checkout:
#
#   sudo ./dune install-service
#
# Manual install:
#   1. Replace @DUNE_USER@, @DUNE_GROUP@, and @PROJECT_ROOT@ below.
#   2. sudo cp scripts/dune-stack.service /etc/systemd/system/dune-stack.service
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now dune-stack.service
#
# `dune start`/`dune stop` already wrap docker compose with the right profile
# and overlays (DEPLOYMENT_PROFILE, DUNE_HOSTNET_OVERLAY), so this unit just
# invokes the CLI. Pair it with `sudo ./scripts/cpu-pin.sh --install` to also
# restore CPU pinning on boot.

[Unit]
Description=Dune Awakening dedicated server stack (docker compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=@DUNE_USER@
Group=@DUNE_GROUP@
WorkingDirectory=@PROJECT_ROOT@
ExecStart=@PROJECT_ROOT@/dune start
ExecStop=@PROJECT_ROOT@/dune stop
# Container bring-up can take a while on first boot (image loads, DB init).
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
